# community_scripts/karakeep-time-tagger/karakeep-time-tagger.service
[Unit]
Description=Add time-to-read tags to Karakeep bookmarks by calling Karakeep-Time-Tagger
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
# Update this path to your karakeep-python-api repository location
WorkingDirectory=%h/repos/karakeep-python-api
ExecStart=/usr/bin/python3 community_scripts/karakeep-time-tagger/karakeep-time-tagger.py --verbose=false
# Uncomment and set path to environment file containing KARAKEEP_PYTHON_API_KEY and KARAKEEP_PYTHON_API_BASE_URL
EnvironmentFile=%h/.config/karakeep/env
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
